# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tomcat9]
command=/usr/share/tomcat9/bin/catalina.sh run
directory=/usr/share/tomcat9
user=tomcat
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/tomcat9.err.log
stdout_logfile=/var/log/supervisor/tomcat9.out.log
environment=JAVA_HOME="/usr/lib/jvm/default-java",CATALINA_HOME="/usr/share/tomcat9",CATALINA_BASE="/var/lib/tomcat9",CATALINA_TMPDIR="/tmp"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface 